# src/xwared.service.template
[Unit]
Description=Xware Service
After=network.target

[Service]
Type=simple
ExecStart=##PREFIX##/xwared

[Install]
WantedBy=default.target